# src/files/etc/rc.local
/usr/local/share/isotop/bin/zerohosts /etc/hosts.include &
syspatch &
