# src/files/etc/rc.local
days=7
if [ -n "$(find /var/unwind.block -mtime +${days} )" ]; then
    /usr/local/sbin/unwind-block > /var/unwind.block &
fi
syspatch &
